# deploy/supervisor/breast-ai-frontend.conf
[program:breast-ai-frontend]
command=/usr/bin/npm run preview -- --host 0.0.0.0 --port 3000
directory=/workspace/breast-ai-system/frontend
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/breast-ai/frontend.log
stdout_logfile_maxbytes=50MB
environment=NODE_ENV="production"
