# deploy/supervisor/breast-ai-backend.conf
[program:breast-ai-backend]
command=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
directory=/workspace/breast-ai-system/backend
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/breast-ai/backend.log
stdout_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1",DATABASE_URL="postgresql://breast_ai_user:secure_password_2026@localhost:5432/breast_ai_db"
